# pkg/tproxy-manager-ipk/etc/uci-defaults/90_tproxy_manager
#!/bin/sh
# Idempotent defaults: создаём только отсутствующие секции/параметры.
set -e

# --- /etc/xray: создать директорию и пустые файлы ТОЛЬКО если их нет ---
mkdir -p /etc/xray
for f in \
  /etc/xray/xray-tproxy.ports \
  /etc/xray/xray-tproxy.v4 \
  /etc/xray/xray-tproxy.v6 \
  /etc/xray/xray-tproxy.src4.only \
  /etc/xray/xray-tproxy.src6.only \
  /etc/xray/xray-tproxy.src4.bypass \
  /etc/xray/xray-tproxy.src6.bypass \
  /etc/xray/geo-sources.conf
do
  [ -e "$f" ] || : > "$f"
done

# --- UCI helpers ---
# ensure_section <cfg> <section_name> <type>
ensure_section() {
  local cfg="$1" sec="$2" typ="$3"
  if ! uci -q get "$cfg.$sec" >/dev/null 2>&1; then
    uci set "$cfg.$sec=$typ"
    CHANGED=1
  fi
}
# set_default <cfg> <section> <option> <value>
set_default() {
  local cfg="$1" sec="$2" opt="$3" val="$4" cur
  cur="$(uci -q get "$cfg.$sec.$opt" 2>/dev/null || true)"
  if [ -z "$cur" ]; then
    uci set "$cfg.$sec.$opt=$val"
    CHANGED=1
  fi
}

# --- xray (ядро): только заполняем недостающие ---
CHANGED=0
ensure_section xray enabled enabled
set_default   xray enabled enabled '1'

ensure_section xray config  config
set_default   xray config  confdir '/etc/xray'
set_default   xray config  datadir '/usr/share/xray/'
set_default   xray config  format  'json'
uci -q delete xray.config.conffiles || true
[ "$CHANGED" -eq 1 ] && uci commit xray

# --- xray-proxy (наш UCI): создаём секцию и каждое поле, если отсутствует ---
CHANGED=0
ensure_section xray-proxy main main

set_default xray-proxy main log_enabled        '1'
set_default xray-proxy main nft_table          'xray'
set_default xray-proxy main ifaces             'br-lan'
set_default xray-proxy main ipv6_enabled       '1'
set_default xray-proxy main fwmark_tcp         '0x1'
set_default xray-proxy main fwmark_udp         '0x2'
set_default xray-proxy main rttab_tcp          '100'
set_default xray-proxy main rttab_udp          '101'
set_default xray-proxy main port_mode          'bypass'
set_default xray-proxy main ports_file         '/etc/xray/xray-tproxy.ports'
set_default xray-proxy main bypass_v4_file     '/etc/xray/xray-tproxy.v4'
set_default xray-proxy main bypass_v6_file     '/etc/xray/xray-tproxy.v6'
set_default xray-proxy main src_mode           'off'
set_default xray-proxy main src_only_v4_file   '/etc/xray/xray-tproxy.src4.only'
set_default xray-proxy main src_only_v6_file   '/etc/xray/xray-tproxy.src6.only'
set_default xray-proxy main src_bypass_v4_file '/etc/xray/xray-tproxy.src4.bypass'
set_default xray-proxy main src_bypass_v6_file '/etc/xray/xray-tproxy.src6.bypass'

[ "$CHANGED" -eq 1 ] && uci commit xray-proxy

exit 0