# pkg/tproxy-manager-ipk/etc/uci-defaults/90_tproxy_manager
#!/bin/sh
# Idempotent defaults for OpenWrt UCI

set -e

# --- /etc/tproxy-manager: каталог и пустые файлы только если их нет ---
mkdir -p /etc/tproxy-manager
for f in \
  /etc/tproxy-manager/tproxy-manager.ports \
  /etc/tproxy-manager/tproxy-manager.v4 \
  /etc/tproxy-manager/tproxy-manager.v6 \
  /etc/tproxy-manager/tproxy-manager.src4.only \
  /etc/tproxy-manager/tproxy-manager.src6.only \
  /etc/tproxy-manager/tproxy-manager.src4.bypass \
  /etc/tproxy-manager/tproxy-manager.src6.bypass \
  /etc/tproxy-manager/geo-sources.conf
do
  [ -e "$f" ] || : > "$f"
done

# --- UCI доступен? ---
command -v uci >/dev/null 2>&1 || exit 0

# На случай полностью отсутствующих файлов конфигов (безопасно):
for c in tproxy-manager tproxy-manager; do
  [ -f "/etc/config/$c" ] || touch "/etc/config/$c"
done

# --- UCI helpers ---
# set_default <cfg> <section> <option> <value>
set_default() {
  local cfg="$1" sec="$2" opt="$3" val="$4"
  if ! uci -q get "$cfg.$sec.$opt" >/dev/null 2>&1; then
    uci set "$cfg.$sec.$opt=$val"
  fi
}

# ensure_section <cfg> <section_name> <type>
ensure_section() {
  local cfg="$1" sec="$2" typ="$3"
  # Универсально: создаёт пакет и именованную секцию при отсутствии
  uci set "$cfg.$sec=$typ"
}

# --- tproxy-manager (наш UCI) ---
ensure_section tproxy-manager main main
set_default tproxy-manager main log_enabled        '1'
set_default tproxy-manager main nft_table          'tproxy'
set_default tproxy-manager main ifaces             'br-lan'
set_default tproxy-manager main ipv6_enabled       '1'
set_default tproxy-manager main tproxy_port        '61219'
set_default tproxy-manager main fwmark_tcp         '0x1'
set_default tproxy-manager main fwmark_udp         '0x2'
set_default tproxy-manager main rttab_tcp          '100'
set_default tproxy-manager main rttab_udp          '101'
set_default tproxy-manager main port_mode          'bypass'
set_default tproxy-manager main ports_file         '/etc/tproxy-manager/tproxy-manager.ports'
set_default tproxy-manager main bypass_v4_file     '/etc/tproxy-manager/tproxy-manager.v4'
set_default tproxy-manager main bypass_v6_file     '/etc/tproxy-manager/tproxy-manager.v6'
set_default tproxy-manager main src_mode           'off'
set_default tproxy-manager main src_only_v4_file   '/etc/tproxy-manager/tproxy-manager.src4.only'
set_default tproxy-manager main src_only_v6_file   '/etc/tproxy-manager/tproxy-manager.src6.only'
set_default tproxy-manager main src_bypass_v4_file '/etc/tproxy-manager/tproxy-manager.src4.bypass'
set_default tproxy-manager main src_bypass_v6_file '/etc/tproxy-manager/tproxy-manager.src6.bypass'
uci commit tproxy-manager

exit 0